# .config/etc/udev/rules.d/95-monitor.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/lukas/.Xauthority", RUN+="/home/lukas/.scripts/hotplug-monitor.sh"
